# .docker/hello-world/etc/service/apache2/run
#!/bin/sh

exec httpd -D FOREGROUND -f /etc/apache2/httpd.conf